# docker/supervisord.conf
[supervisord]
nodaemon=true

[program:redis-server]
command=/usr/bin/redis-server

[program:nodepot]
command=/bin/bash -c "cd /opt/Nodepot/ && nodejs app.js"